# teamviewerd.service.d/override.conf
[Service]
KillSignal=SIGKILL
TimeoutStopSec=5